# meta-cpane/recipes-core/cpane-users/files/set-cpane-password.service
[Unit]
Description=Set cpane user password on first boot
After=multi-user.target
ConditionPathExists=!/var/lib/cpane-password-set

[Service]
Type=oneshot
ExecStart=/bin/sh -c 'echo "cpane:cpane" | chpasswd && touch /var/lib/cpane-password-set'
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target

